# 60-usb_nfc.rules
ACTION=="add|remove",SUBSYSTEM=="usb",ENV{ID_VENDOR_ID}=="1db2",ENV{ID_MODEL_ID}=="060f",RUN+="/usr/local/bin/usb-nfc-$env{ACTION}.sh"
